# claude-proxy/deploy/claude-proxy.service
[Unit]
Description=Claude Proxy Service - Claude CLI Gateway
Documentation=https://github.com/r9r-dev/home-agent
After=network.target
Wants=network-online.target

[Service]
Type=simple

# Run as the user who has Claude CLI configured
# Change this to your username
User=rlamour
Group=staff

# Working directory
WorkingDirectory=/opt/claude-proxy

# Binary path
ExecStart=/opt/claude-proxy/claude-proxy

# Environment variables
# Port and host configuration
Environment="PROXY_PORT=9090"
Environment="PROXY_HOST=0.0.0.0"

# Path to Claude CLI binary (adjust if not in PATH)
Environment="CLAUDE_BIN=claude"

# API Key for authentication (generate a secure random key)
# Example: openssl rand -hex 32
Environment="PROXY_API_KEY="

# Restart policy
Restart=always
RestartSec=5s

# Process limits
LimitNOFILE=65535

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-proxy

[Install]
WantedBy=multi-user.target
